# scripts/com.claus.dailystack.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claus.dailystack</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/clauseduardpetraeus/Github repos/daily-stack/.venv/bin/uvicorn</string>
        <string>server.app:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8080</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/clauseduardpetraeus/Github repos/daily-stack</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/dailystack.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/dailystack.err</string>
</dict>
</plist>
